# packages/gancho/gancho-0.1.0.tar.gz/gancho-0.1.0/gancho.socket
[Unit]
Description=gancho socket

[Socket]
ListenStream=/run/gancho.sock

[Install]
WantedBy=sockets.target
