# packages/catsoop/catsoop-19.0.7.tar.gz/catsoop-19.0.7/systemd/catsoop.service
[Unit]
Description=CAT-SOOP
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/catsoop start
User=catsoop
KillMode=mixed
TimeoutSec=10

[Install]
WantedBy=multi-user.target
