# packages/dlrepo/dlrepo-0.44.tar.gz/dlrepo-0.44/etc/dlrepo.socket
# ex: ft=systemd
[Unit]
Description=Socket for dlrepo artifact server

[Socket]
ListenStream=/run/dlrepo.sock
SocketUser=dlrepo
SocketGroup=www-data
SocketMode=0660

[Install]
WantedBy=sockets.target
