# packages/dlrepo/dlrepo-0.44.tar.gz/dlrepo-0.44/tests/publish_reference/branches/branch/tag/job1/fmt2/dlrepo-cleanup.timer
# ex: ft=systemd

[Unit]
Description=Cleanup tags according to branch cleanup policy.

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
