# packages/swapdog/swapdog-0.2.2.tar.gz/swapdog-0.2.2/swapdog.service
[Unit]
Description=SwapDog Service - Enables swaps at set thresholds
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/sbin/swapdog.py /etc/swapdog.json
TimeoutSec=10
KillMode=process
Restart=on-failure
RestartSec=5
User=root
Group=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
